# build/sysroot/usr/lib/systemd/system/casaos.service
[Unit]
After=casaos-message-bus.service
After=rclone.service
Description=CasaOS Main Service

[Service]
ExecStart=/usr/bin/casaos -c /etc/casaos/casaos.conf
PIDFile=/var/run/casaos/casaos.pid
Restart=always
Type=notify

[Install]
WantedBy=multi-user.target
